# bbs.conf
[program:bbs]
command=/root/.virtualenvs/xjmweb/bin/gunicorn wsgi -c gunicorn.config.py
directory=/root/WebProjects/webbbs
autostart=true
autorestart=true
