# bbs.conf
[program:bbs]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/root/WebProjects/webbbs
autostart=true
autorestart=true
